# scripts/configs/.local_bin/arrpc.service
[Service]
Type=simple
ExecStart=arrpc
WorkingDirectory=~

[Install]
WantedBy=default.target
